# s6-overlay/s6-rc.d/code-server/run
#!/command/execlineb -P
code-server --host 0.0.0.0 --auth none --port 8080