# s6-overlay/s6-rc.d/jupyter-lab/run
#!/command/execlineb -P
jupyter lab --LabApp.base_url=jl --ip 0.0.0.0 --allow-root --LabApp.token="" --LabApp.password="" --LabApp.open_browser="False"